# com.shft.screensaver.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shft.screensaver.cleanup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>while true; do pgrep -x legacyScreenSaver &gt;/dev/null 2&gt;&amp;1 &amp;&amp; ! pgrep -x ScreenSaverEngine &gt;/dev/null 2&gt;&amp;1 &amp;&amp; killall legacyScreenSaver 2&gt;/dev/null; sleep 10; done</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
